# chrome-container/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=5
loglevel=info
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# Ensure log directory exists
[program:prepare]
command=bash -c "mkdir -p /var/log && chmod -R 777 /var/log"
startsecs=0
autorestart=false
startretries=1
priority=1

[program:xvfb]
command=/app/start-xvfb.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err
priority=1

[program:x11vnc]
command=/app/start-vnc.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err
priority=2
depends_on=xvfb
startretries=5

[program:fluxbox]
command=/usr/bin/fluxbox
autorestart=true
stdout_logfile=/var/log/fluxbox.log
stderr_logfile=/var/log/fluxbox.err
priority=30
depends_on=xvfb
startretries=5

[program:novnc]
command=/app/start-novnc.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err
priority=3
depends_on=x11vnc
startretries=5

[program:chrome]
command=/app/start-chrome.sh
autostart=true
autorestart=true
startsecs=10
startretries=3
user=chrome
environment=HOME="/home/chrome",USER="chrome"
stdout_logfile=/var/log/chrome.log
stderr_logfile=/var/log/chrome.err
priority=100
depends_on=xvfb,x11vnc,fluxbox 